# msys2/usr/share/doc/gnupg/examples/systemd-user/keyboxd.service
[Unit]
Description=GnuPG public key management service
Documentation=man:keyboxd(8)
Requires=keyboxd.socket

[Service]
ExecStart=/usr/lib/gnupg/keyboxd --supervised
ExecReload=/usr/bin/gpgconf --reload keyboxd
