# msys2/var/lib/pacman/local/pacman-contrib-1.10.6-1/files
%FILES%
usr/
usr/bin/
usr/bin/checkupdates
usr/bin/paccache
usr/bin/pacdiff
usr/bin/paclist
usr/bin/paclog-pkglist
usr/bin/pacscripts
usr/bin/pacsearch
usr/bin/pacsort.exe
usr/bin/pactree.exe
usr/bin/rankmirrors
usr/bin/updpkgsums
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/paccache.service
usr/lib/systemd/system/paccache.timer
usr/lib/systemd/system/pacman-filesdb-refresh.service
usr/lib/systemd/system/pacman-filesdb-refresh.timer
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/checkupdates.8.gz
usr/share/man/man8/paccache.8.gz
usr/share/man/man8/pacdiff.8.gz
usr/share/man/man8/paclist.8.gz
usr/share/man/man8/paclog-pkglist.8.gz
usr/share/man/man8/pacscripts.8.gz
usr/share/man/man8/pacsearch.8.gz
usr/share/man/man8/pacsort.8.gz
usr/share/man/man8/pactree.8.gz
usr/share/man/man8/rankmirrors.8.gz
usr/share/man/man8/updpkgsums.8.gz
usr/share/vim/
usr/share/vim/vimfiles/
usr/share/vim/vimfiles/ftdetect/
usr/share/vim/vimfiles/ftdetect/PKGBUILD.vim
usr/share/vim/vimfiles/ftdetect/SRCINFO.vim
usr/share/vim/vimfiles/ftplugin/
usr/share/vim/vimfiles/ftplugin/PKGBUILD.vim
usr/share/vim/vimfiles/syntax/
usr/share/vim/vimfiles/syntax/PKGBUILD.vim
usr/share/vim/vimfiles/syntax/SRCINFO.vim
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_checkupdates
usr/share/zsh/site-functions/_paccache
usr/share/zsh/site-functions/_pacdiff
usr/share/zsh/site-functions/_paclist
usr/share/zsh/site-functions/_paclog-pkglist
usr/share/zsh/site-functions/_pacscripts
usr/share/zsh/site-functions/_pacsearch
usr/share/zsh/site-functions/_pacsort
usr/share/zsh/site-functions/_pactree
usr/share/zsh/site-functions/_rankmirrors
usr/share/zsh/site-functions/_updpkgsums

