# msys2/var/lib/pacman/local/pacman-contrib-1.3.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/checkupdates
usr/bin/paccache
usr/bin/pacdiff
usr/bin/paclist
usr/bin/paclog-pkglist
usr/bin/pacscripts
usr/bin/pacsearch
usr/bin/pacsort.exe
usr/bin/pactree.exe
usr/bin/rankmirrors
usr/bin/updpkgsums
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/paccache.service
usr/lib/systemd/system/paccache.timer
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/checkupdates.8.gz
usr/share/man/man8/paccache.8.gz
usr/share/man/man8/pacdiff.8.gz
usr/share/man/man8/pacsort.8.gz
usr/share/man/man8/pactree.8.gz
usr/share/man/man8/updpkgsums.8.gz
usr/share/vim/
usr/share/vim/vimfiles/
usr/share/vim/vimfiles/ftdetect/
usr/share/vim/vimfiles/ftdetect/PKGBUILD.vim
usr/share/vim/vimfiles/ftplugin/
usr/share/vim/vimfiles/ftplugin/PKGBUILD.vim
usr/share/vim/vimfiles/syntax/
usr/share/vim/vimfiles/syntax/PKGBUILD.vim

