# msys2/usr/lib/systemd/system/paccache.service
[Unit]
Description=Remove unused cached package files

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
