# msys2/usr/lib/systemd/system/paccache.service
[Unit]
Description=Discard unused packages

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
